# systemd/codex-telegram-worker.service
[Unit]
Description=Codex Telegram Worker
After=network-online.target

[Service]
Type=simple
UMask=0077
WorkingDirectory=%h/codex-telegram-bot
EnvironmentFile=%h/codex-telegram-bot/.env
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env node %h/codex-telegram-bot/src/worker.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
